# php/5.6.40/etc/supervisord.d/nginx.ini
[program:nginx]
command = /usr/local/nginx/sbin/nginx
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE
